# oc-set-root-sshkeys.service
[Unit]
Description=Setting root SSH Keys

[Service]
Type=oneshot
ExecStart=@libexec@/set-root-sshkeys

[Install]
WantedBy=basic.target
